# arch_challenges/kippo/systemd/kippo.service
[Unit]
Description=Kippo Honeypot

[Service]
User=kippo
Group=kippo
ExecStart=/home/kippo/startkippo

[Install]
WantedBy=multi-user.target
